# sample.socket
[Unit]
Description=Web Application Socket

[Socket]
ListenStream=127.0.0.1:8888

[Install]
WantedBy=sockets.target
